# var/lib/dpkg/info/chrony.conffiles
/etc/NetworkManager/dispatcher.d/20-chrony
/etc/default/chrony
/etc/init.d/chrony
/etc/logrotate.d/chrony
/etc/network/if-post-down.d/chrony
/etc/network/if-up.d/chrony
/etc/ppp/ip-down.d/chrony
/etc/ppp/ip-up.d/chrony
