# var/lib/dpkg/info/chrony.list
/.
/etc
/etc/NetworkManager
/etc/NetworkManager/dispatcher.d
/etc/NetworkManager/dispatcher.d/20-chrony
/etc/chrony
/etc/default
/etc/default/chrony
/etc/init.d
/etc/init.d/chrony
/etc/logrotate.d
/etc/logrotate.d/chrony
/etc/network
/etc/network/if-post-down.d
/etc/network/if-post-down.d/chrony
/etc/network/if-up.d
/etc/network/if-up.d/chrony
/etc/ppp
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/chrony
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/chrony
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/chrony.service
/usr
/usr/bin
/usr/bin/chronyc
/usr/sbin
/usr/sbin/chronyd
/usr/share
/usr/share/chrony
/usr/share/chrony/chrony.conf
/usr/share/chrony/chrony.keys
/usr/share/doc
/usr/share/doc/chrony
/usr/share/doc/chrony/FAQ.gz
/usr/share/doc/chrony/NEWS.Debian.gz
/usr/share/doc/chrony/README.Debian
/usr/share/doc/chrony/README.gz
/usr/share/doc/chrony/changelog.Debian.gz
/usr/share/doc/chrony/changelog.gz
/usr/share/doc/chrony/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/chrony
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/chronyc.1.gz
/usr/share/man/man5
/usr/share/man/man5/chrony.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/chronyd.8.gz
/usr/share/menu
/usr/share/menu/chrony
/var
/var/lib
/var/lib/chrony
/var/log
/var/log/chrony
